# config/com.user.recipe-recommender.send.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.recipe-recommender.send</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yeesheen/projects/recipe-recommender/venv/bin/python</string>
        <string>-m</string>
        <string>src.scheduler.send_daily</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yeesheen/projects/recipe-recommender</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/yeesheen/projects/recipe-recommender/logs/daily_send.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yeesheen/projects/recipe-recommender/logs/daily_send_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
